# home/user/streamdeck/streamdeck.service
[Unit]
Description=Streamdeck Python Service


[Service]
ExecStart=/home/user/streamdeck/script.sh

